# lib/flavors/voidlinux/files/etc/cron.hourly/btrbk
#!/bin/sh
exec /usr/sbin/btrbk -q run > /dev/null 2>&1
