# lib/flavors/voidlinux/files/etc/cron.monthly/btrfs-srub
#!/bin/sh
exec /bin/btrfs scrub start / > /dev/null 2>&1
